# server/etc/systemd/system/certbot.service
[Unit]
Description=Let's Encrypt renewal

[Service]
ExecStart=/opt/certbot/bin/certbot renew -q --post-hook "systemctl reload nginx.service"
